# Taxi/mobile_qa_test_app/mobile-qa-app-supervisor.conf
[program:mobile-qa-test-app]
command = su -c "dart run conduit:conduit serve --port 80"
directory = /app/backend
autostart = yes
autorestart = yes
startretries = 6
stdout_logfile_maxbytes = 50MB
stderr_logfile_maxbytes = 50MB
stdout_logfile_backups = 2
stderr_logfile_backups = 2
stdout_logfile = /var/log/supervisor/mobile-qa-test-app.log
stderr_logfile = /var/log/supervisor/mobile-qa-test-app.log


