# Taxi/testing/celery_task_status_worker.conf
[program:celery_task_status_worker]
environment=PYTHONPATH="/usr/lib/yandex/startek-graph:/usr/lib/yandex/startek-graph/api",USER="root"
command=/usr/lib/yandex/startek-graph/bin/celery worker -A api.workers.task_status_worker --loglevel=INFO --concurrency=2 -Q startrek_graph_task_status
directory=/usr/lib/yandex/startek-graph/

user=root

numprocs=1
process_name=%(program_name)s

stdout_logfile=/var/log/yandex/startrek_graph/celery_task_status_worker.log
stderr_logfile=/var/log/yandex/startrek_graph/celery_task_status_worker_error.log

autostart=true
autorestart=true
startsecs=0
stopwaitsecs=30
killasgroup=true
stopasgroup=true