# Taxi/testing/celery_worker_test_backup_big.conf
[program:yandex_taxi_dmp_gpdb_manager_celery_worker_test_backup_big]
environment=PYTHONPATH="/usr/lib/yandex/gpdb-manager",USER="root"

command=/usr/lib/yandex/gpdb-manager/.env/bin/celery worker -A workers.backup --loglevel=INFO --concurrency=3 -Q test_backup_big

directory=/usr/lib/yandex/gpdb-manager/

user=root

numprocs=1
process_name=%(program_name)s

stdout_logfile=/var/log/yandex/taxidwh-gpdb-manager/worker-test-backup-big.log
stderr_logfile=/var/log/yandex/taxidwh-gpdb-manager/worker-test-backup-big-error.log

autostart=true
autorestart=true
startsecs=0
stopwaitsecs=30
killasgroup=true
stopasgroup=true