# Taxi/yandex-mobile-test-taxi-tracker-mongodb-config/src/taxi-mongo-tmpfs.conf
# taxi-mongo-tmpfs - mounts /opt/mongodb in tmpfs before mongodb starts
#
description     "Mount /opt/mongodb/ on boot"

start on (startup or (mounted))
stop on runlevel [06]

task

env MONGOMOUNT=/opt/mongodb/

script

        if ! mountpoint -q "$MONGOMOUNT"
        then
                mount -t tmpfs -o uid=mongodb,gid=mongodb,size=4096m tmpfs "$MONGOMOUNT"
        fi
end script

pre-stop script
umount "$MONGOMOUNT"
end script

