# Taxi/GENERAL/suptech_bot_testing.conf
[program:suptech_bot]
directory=/app/suptech-bot
command=/usr/lib/yandex/taxi-py3-2/bin/python /app/suptech-bot/suptech_bot/catbot_app.py
autostart=true
autorestart=true
stdout_logfile=/var/log/catbot_out_main
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/catbot_error_main
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/usr/lib/yandex/taxi-scripts-cron"

[program:telegram_notifator]
command=/usr/lib/yandex/taxi-py3-2/bin/python /app/suptech-bot/other_bots/telegram_notifator.py
stdout_logfile=/var/log/telegram_notifator_out.log
stderr_logfile=/var/log/telegram_notifator_error.log
autostart=true
autorestart=true
environment=PYTHONPATH="/usr/lib/yandex/taxi-scripts-cron"

[program:magistral_bot]
command=/usr/lib/yandex/taxi-py3-2/bin/python /app/suptech-bot/other_bots/magistral_bot/app.py
autostart=true
autorestart=true
stdout_logfile=/var/log/magistral_bot_out_main
stderr_logfile=/var/log/magistral_bot_error_main
environment=PYTHONPATH="/usr/lib/yandex/taxi-scripts-cron"
