# Taxi/GENERAL/yandex-taxi-octocore.testing (2).conf
[program:taxi-octocore]
process_name = %(program_name)s_%(process_num)02d
command = /usr/bin/octocore.py 80%(process_num)02d
user = www-data
group = www-data
numprocs = 4
numprocs_start = 1
autostart = yes
autorestart = yes
startsecs = 5
stopsignal = TERM
startretries = 10
stopwaitsecs = 10
stdout_logfile = /var/log/supervisor/taxi-octocore.log
stderr_logfile = /var/log/supervisor/taxi-octocore.log
