# Taxi/test_rtc_etl/workflow/crontabs/workflow-cron-production
# chkcrontab: disable-msg=HOURS_NOT_MINUTES

MAILTO="taxi-dwh-cron@yandex-team.ru"

PATH="/usr/lib/yandex/taxi-dmp-deps-py3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
PYTHONPATH=/usr/lib/yandex/taxi-dmp-test-rtc-etl
LANG=en_US.UTF-8
LRPY="python3.7 -m dmp_suite.runner python"
LRSH="python3.7 -m dmp_suite.runner sh"

# ##############################################################################
# ##############################  PROCESS KILLER  ####################################
 */2    * * * *    www-data       $LRPY --no-scheduler-lock -L "dmp_suite.maintenance.process_killer.cmd"
 10    22 * * *    www-data       $LRPY --no-scheduler-lock -L "dmp_suite.maintenance.process_killer.cmd_clean"
