# Taxi/test_memlockd/build/services/test-service/memlock.conf
[program:memlockd]
command=/usr/sbin/memlockd-test-service
user = www-data
autostart = yes
autorestart = yes
startsecs = 2
stopsignal = TERM
stopwaitsecs = 30
directory = /
stdout_logfile = /var/log/supervisor/memlockd.stdout.log
stderr_logfile = /var/log/supervisor/memlockd.stderr.log
