# Taxi/test_supervisor/test_full/test-service/debian/test-package-name.conf
[program:simple-program]
process_name = %(program_name)s_%(process_num)02d
numprocs = 30
command = sleep infinity
user = alberist
umask = 022
numprocs_start = 10
autostart = yes
autorestart = yes
startsecs = 15
stopsignal = TERM
startretries = 20
stopwaitsecs = 25
directory = /home/me/my/service
environment = PATH="$PATH:/usr/lib/uber/taxi-deps-py3-2/bin",PYTHONPATH="/usr/lib/uber/taxi-stq"
stdout_logfile = /var/log/supervisor/simple-program.log
stderr_logfile = /var/log/supervisor/simple-program.log
killasgroup = true
