# Taxi/test_supervisor/test_simple/test-service/debian/test-package-name.conf
[program:simple-program]
command = sleep infinity
user = www-data
umask = 022
autostart = yes
autorestart = yes
stopsignal = TERM
directory = /home/me/my/service
stdout_logfile = /var/log/supervisor/simple-program.log
stderr_logfile = /var/log/supervisor/simple-program.log
